# yocto_training_materials/code/meta-allegion/meta-allegion/recipes-example/mydaemon/files/mydaemon.service
[Unit]
Description=Allegion example daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/mydaemon
Restart=on-failure
RestartSec=5

# Run as the unprivileged user created by the recipe
User=mydaemon
Group=mydaemon

# --- Hardening (see Linux Capabilities study guide) ---
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6
RestrictNamespaces=yes
CapabilityBoundingSet=
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# --- Resource limits (see cgroups study guide) ---
MemoryMax=50M
TasksMax=16
CPUQuota=20%

[Install]
WantedBy=multi-user.target
